# github-auto-pr.service
[Unit]
Description=GitHub Auto PR Bot Service
After=network.target

[Service]
Type=oneshot
User=your_username
WorkingDirectory=/home/your_username/apps/auto-activity
Environment=GITHUB_TOKEN=your_github_token_here
Environment=GITHUB_USERNAME=your_github_username
Environment=GITHUB_REPO=your_repo_name
ExecStart=/home/your_username/apps/auto-activity/venv/bin/python /home/your_username/apps/auto-activity/github-auto-pr.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
